# blog-sanic.conf
[program:blog-sanic]
command=/usr/local/bin/gunicorn wsgi:application -c conf-gunicorn.py
directory=/var/www/blog-sanic
autostart=true
autorestart=true